# autostart/zzz-conky-ddd.desktop
[Desktop Entry]
Exec=conky
Terminal=false
Type=Application
OnlyShowIn=OPENBOX;
