# autostart/zzz-volumeicon-ddd.desktop
[Desktop Entry]
Exec=volumeicon
Terminal=false
Type=Application
OnlyShowIn=OPENBOX;
